# mytimer.timer
[Unit]
Description=Runs mytimer.service every 10s

[Timer]
OnBootSec=1s
OnUnitActiveSec=10s
Unit=mytimer.service

[Install]
WantedBy=timers.target

[Timer]

[Install]
